# 15-mouse.rules
ACTION==”remove”, ID_BUS==”usb??”, RUN+="sh /home/will/Desktop/disable-mouse-in-linux/script.sh"