# 50-mouse.rules
SUBSYSTEM=="usb_device",ACTION==”add”,RUN+="/usr/bin/disable-mouse"