# disable-mouse.service
[Unit]
Description=Disable Mouse
[Service]
ExecStart=/bin/bash /bin/disable-mouse/disable-mouse.sh
[Install]
WantedBy=multi-user.target