# disable-mouse.service
[Unit]
Description=Disable Mouse

[Service]
ExecStart=/bin/python3 /bin/disable-mouse/automated.py